# supervisord.conf
[supervisord]
nodaemon=true

[program:flask]
command=python /data/manage.py runserver -h 0.0.0.0 -p 80
